# development/gcc13/profile.d/gcc13.csh
setenv CC     gcc-13
setenv CPP    cpp-13
setenv CXX    g++-13
setenv AR     gcc-ar-13
setenv NM     gcc-nm-13
setenv RANLIB gcc-ranlib-13

