# development/gcc13/profile.d/gcc13.sh
export CC=gcc-13
export CPP=cpp-13
export CXX=g++-13
export AR=gcc-ar-13
export NM=gcc-nm-13
export RANLIB=gcc-ranlib-13

